# scripts/monit/centos7/files/crontab/alive
*/2 * * * * pzoo /home/spanel/modules/pastafari/scripts/monit/centos7/files/get_info.py

