# scripts/monit/debian_jessie/files/crontab/alive
*/2 * * * * pzoo /home/spanel/vendor/phangoapp/leviathan/scripts/monit/debian_jessie/files/get_info.py

